# redbot-restart@.service
[Unit]
Description=restart redbot under podman running from %i

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart redbot@%i
